# config/crontab
0 0 * * * python3 /path/to/your/project/auto/main.py